# monitorix/supervisord.conf
[supervisord]
nodaemon=true

[program:monitorix]
command=/usr/bin/monitorix -c /etc/monitorix/monitorix.conf -p /run/monitorix.pid

